# crontab
SHELL=/bin/bash
0 22 * * * /app/rebuild.py >> /dev/pts/0 2>&1
