# crontab
SHELL=/bin/bash
0 22 * * * /app/rebuild.py >> /build/log
